# cache-weather.service
[Unit]
Description=Cache Weather Data

[Service]
ExecStart=/usr/bin/fetch-weather.sh

